# kafka_zookeeper/docker-compose/supervisord.conf
[supervisord]
nodaemon=true

[program:service]
nodaemon=true
command=/usr/bin/start-service.sh
