# kafka_zookeeper/docker-compose/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
nodaemon=true
autostart=true
autorestart=true
